# systemd/mpd-notification.service
[Unit]
Description=MPD Notification

[Service]
Type=simple
ExecStart=/usr/bin/mpd-notification

[Install]
WantedBy=default.target
